# deploy/systemd/congopro-bridge.service
[Unit]
Description=Congopro Bridge
After=network.target
# Add "After=traefik.service" if you want strict ordering.

[Service]
Type=simple
User=ops
Group=ops

Environment=PORT=8110
Environment=OLLAMA_URL=http://localhost:11434
Environment=GENERATIVE_MODEL=gemma3:1b
Environment=EMBEDDING_MODEL=nomic-embed-text
Environment=MODELS_DIR=/opt/congopro-bridge/models
Environment=ALLOWED_ORIGIN=https://congopro.com

WorkingDirectory=/opt/congopro-bridge
ExecStart=/opt/congopro-bridge/congopro-bridge
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=congopro-bridge

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/congopro-bridge
ProtectHome=true

[Install]
WantedBy=multi-user.target