# deploy/systemd/congopro-bridge.service
[Unit]
Description=Congopro Bridge
After=network.target
# Add "After=traefik.service" if you want strict ordering.

[Service]
Type=simple
User=ops
Group=ops

Environment=PORT=8110
Environment=OLLAMA_URL=http://127.0.0.1:11434
Environment=GENERATIVE_MODEL=gemma3:1b
Environment=EMBEDDING_MODEL=nomic-embed-text
Environment=ALLOWED_ORIGIN=https://congopro.com
Environment=OLLAMA_ALLOW_PUBLIC_IP=false
Environment=OLLAMA_ALLOWED_HOSTS=
Environment=MEILI_URL=http://127.0.0.1:7700
Environment=MEILI_MASTER_KEY=slcliuSNfGvYciwfoSFejvk9vVYd5xIL
Environment=MEILI_INDEX_NAME=companies

WorkingDirectory=/opt/congopro-bridge
ExecStart=/opt/congopro-bridge/congopro-bridge
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=congopro-bridge

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/congopro-bridge
ProtectHome=true

[Install]
WantedBy=multi-user.target